# etc/systemd/system/sds011.service
[Unit]
Description=SDS011 Air Quality Sensor Publisher
After=network.target

[Service]
Type=simple
ExecStart=/home/particles/particles-venv/bin/python /home/particles/airquality/bin/user/sds011.py
WorkingDirectory=/home/particles/airquality/bin/user
Restart=always
RestartSec=10
User=particles
Group=particles
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
